# lib/systemd/system/certcopyToService.service
[Unit]
Description=service to copy certificates to service folders

[Service]
Type=oneshot
ExecStart=/opt/scripts/certcopyToService.sh
