# lib/systemd/system/certcopyFromHost.service
[Unit]
Description=service to copy certificates from host to opt

[Service]
Type=oneshot
ExecStart=/opt/scripts/certcopyFromHost.sh
