# lib/systemd/system/certcopyToOpt.service
[Unit]
Description=service to copy certificates to opt folder

[Service]
Type=oneshot
ExecStart=/opt/scripts/certcopyToOpt.sh
