# lib/systemd/system/cleanup-github-runner.service
[Unit]
Description=service to cleanup github runner

[Service]
Type=oneshot
ExecStart=/opt/scripts/cleanupGithubRunner.sh