# lib/systemd/system/docker-container-restart.service
[Unit]
Description=service to restart rocketchat docker container

[Service]
Type=oneshot
ExecStart=/usr/bin/docker container restart rocketchat