# lxc/pihole/lib/systemd/system/pihole-update-gravity.service
[Unit]
Description=service to update pihole-gravity

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole -g