# lxc/pihole/lib/systemd/system/pihole-updatlists.service
[Unit]
Description=service to update pihole-updatelists

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pihole-updatelists
SuccessExitStatus=0 1