# etc/systemd/system/container-getty@.service.d/override.conf
# this file is used to autologin as root over console (tty) like on proxmox webgui
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear --keep-baud tty%I 115200,38400,9600 $TERM
