# docker/compose/azuracast/restart-azuracast.service
[Unit]
Description=service to restart azuracast 

[Service]
User=docker
ExecStart=/usr/bin/docker container restart azuracast