# docker/lib/systemd/system/certbot-certcopy.service
[Unit]
Description=service to copy certificates to service folders

[Service]
ExecStart=/opt/certbot-certcopy.sh
