# docker/lib/systemd/system/docker-image-prune.service
[Unit]
Description=service to run docker image prune

[Service]
User=docker
ExecStart=/home/docker/bin/docker image prune -f
