# etc/network/if-up.d/vmbr0-ipconfig
#!/bin/sh
[ "$IFACE" = "vmbr0" ] || exit 0
/opt/scripts/network/post-up.sh
